# tests/units/basic.service
[Service]
ExecStart=/bin/true
NoNewPrivileges=true
ProtectSystem=strict
