# tests/units/plus-prefix-info.service
[Service]
ExecStart=/bin/true
RootDirectory=/run/plus-prefix-info-test
ReadOnlyPaths=/etc/example
