# tests/units/plus-prefix-pass.service
[Service]
ExecStart=/bin/true
RootDirectory=/run/plus-prefix-pass-test
ReadOnlyPaths=+/etc/example
