# tests/units/protectproc-no-capfix.service
[Service]
ExecStart=/bin/true
ProtectProc=invisible
