# tests/units/protectproc-capfixed.service
[Service]
ExecStart=/bin/true
ProtectProc=invisible
CapabilityBoundingSet=~CAP_SYS_PTRACE
