# tests/units/protectproc-emptycapset.service
[Service]
ExecStart=/bin/true
ProtectProc=invisible
CapabilityBoundingSet=
