# tests/units/repeated-execpaths.service
[Service]
ExecStart=/bin/true
ExecPaths=/a
ExecPaths=/b
