# examples/privatetmp-noexecpaths-bug.service
# Minimal reproduction of a NoExecPaths= misconfiguration that
# systemd-analyze security cannot detect: NoExecPaths=/tmp below looks
# like it makes PrivateTmp=true's private /tmp non-executable. It does
# not, because RootDirectory= is also set here.
#
# Confirmed by testing (two independent reproductions, 2026-07-24: a
# minimal repro and one matching a real production unit's full hardening
# properties, toggling only RestrictNamespaces=): RootDirectory= alone is
# sufficient to trigger this, with or without RestrictNamespaces= --
# toggling it made no observable difference. Per systemd.exec(5) and
# https://github.com/systemd/systemd/issues/39935: a bare (non-"+") path in
# NoExecPaths=/ExecPaths=/etc. resolves against the host's original root
# instead of RootDirectory= once RootDirectory= is set, silently missing
# the chroot's own private /tmp mount entirely.
#
# Reproduce:
#   sudo systemd-sandbox-check --unit privatetmp-noexecpaths-bug.service
#   -> [FAIL] no_exec_paths (NoExecPaths): ... EXECUTED -- allowlist not enforced
#
# Fix: change NoExecPaths=/tmp to NoExecPaths=+/tmp and re-run:
#   -> [PASS] no_exec_paths (NoExecPaths): ... refused to execute
#
# See systemd.exec(5) and https://github.com/systemd/systemd/issues/39935

[Service]
ExecStart=/usr/bin/sleep infinity
PrivateTmp=true
NoExecPaths=/tmp
RuntimeDirectory=privatetmp-noexecpaths-bug
RootDirectory=/run/privatetmp-noexecpaths-bug
